# DEBIAN/postinst
#!/bin/bash

/usr/sbin/update-rc.d DyalogService defaults

exit 0
